# srcs/requirements/bonus/website/supervisord.conf
[supervisord]
nodaemon=true

[program:myprogram]
command=python /app/myscript.py
autostart=true
autorestart=true
stderr_logfile=/var/log/myprogram.err.log
stdout_logfile=/var/log/myprogram.out.log