# srcs/requirements/bonus/adminer/supervisord.conf
[supervisord]
nodaemon=true

[program:adminer]
command=adminer
autostart=true
autorestart=true
startretries=3